# crontab
# Run maintenance scheduler once daily at 9 AM
0 9 * * * . /tmp/env.sh && cd /app && /usr/local/bin/python /app/src/main.py >> /app/data/cron.log 2>&1

# Empty line at end of file is required for cron
